# roles/configure_monitoring/files/cron-metrics-prometheus.cron
PATH=/usr/sbin:/usr/bin:/sbin:/bin

# read basic metrics
*/15 * * * * root /usr/sbin/univention-metrics-prometheus
