# roles/configure_monitoring/files/cron-metrics-dpx.cron
PATH=/usr/sbin:/usr/bin:/sbin:/bin

# read univention dphoenix specific metrics
*/15 * * * * root /usr/sbin/univention-metrics-dpx
